# infra/launchd/com.her.consolidator.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist template for the her daily consolidator (Phase 4).
  Placeholders __HOME__ and __PROJECT__ are substituted by `bin/her consolidator-install`.

  Schedule: 03:00 daily (CLAUDE.md §5.3).

  ⚠️  IMPORTANT — API keys under launchd:
  direnv is NOT executed by launchd, so GEMINI_API_KEY will NOT be available
  unless you add it explicitly.

  Option A (recommended):
    Add GEMINI_API_KEY to ~/.her/env, then add to EnvironmentVariables below.

  Option B:
    Use `make consolidate` directly from a shell where direnv is active.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.her.consolidator</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>__PROJECT__</string>
        <string>python</string>
        <string>-m</string>
        <string>apps.consolidator</string>
        <string>--db</string>
        <string>__PROJECT__/data/db.sqlite</string>
        <string>--log-dir</string>
        <string>__PROJECT__/data/consolidation_log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT__</string>

    <!-- Run once at 03:00 every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do NOT auto-start on load; only run on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Do not restart on exit -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/.her/logs/consolidator.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.her/logs/consolidator.stderr.log</string>

    <!-- Non-secret environment variables.
         Add GEMINI_API_KEY manually — see header comment. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HER_DATA_DIR</key>
        <string>__PROJECT__/data</string>

        <key>HER_DB_PATH</key>
        <string>__PROJECT__/data/db.sqlite</string>
    </dict>
</dict>
</plist>
